# tts-box/exec-fn-tts-watchdog.timer
[Unit]
Description=Run the exec-fn TTS watchdog every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=15s

[Install]
WantedBy=timers.target
