# api/exec-fn.cron
SHELL=/bin/bash
TZ=America/New_York
30 4 * * * root /app/morning_cron.sh
