# api/exec-fn.cron
SHELL=/bin/bash
30 4 * * * root /app/morning_cron.sh
